# base/service/apache2/run
#!/usr/bin/env bash

exec /usr/sbin/apachectl -DFOREGROUND -k start