# base/service/apache2/run
#!/usr/bin/env bash
exec /usr/sbin/apachectl -DNO_DETACH -k start