# base/service/log-forwarder/run
#!/usr/bin/env bash

exec /etc/service/log-forwarder/log-forwarder.sh