# artifacts/supervisord-cron.conf
[program:cron]
command=/start-cron.sh
numprocs=1
autostart=true
autorestart=true
